# containers/atvloadly/scripts/atvloadly-status-check.service
[Unit]
Description=atvloadly post-refresh status check and notify

[Service]
Type=oneshot
User=hussein
ExecStart=/home/hussein/atvloadly-status-check.sh
